# rc.local
#!/bin/bash

/sbin/ifconfig enp1s0 txqueuelen 65000
exit 0
