# rc.local
#!/bin/bash

/sbin/ifconfig eth0 txqueuelen 65000
exit 0
